# supervisor/dmit_follow.conf
[program:tcat-follow]
command=/usr/bin/php /var/www/html/capture/stream/dmitcat_follow.php
user=www-data
autostart=false
autorestart=true
startretries=30
startsecs=4
