# supervisor/dmit_onepercent.conf
[program:tcat-onepercent]
command=/usr/bin/php /var/www/html/capture/stream/dmitcat_onepercent.php
user=www-data
autostart=false
autorestart=true
startretries=30
startsecs=4
