# supervisor/dmit_track.conf
[program:tcat-track]
command=/usr/bin/php /var/www/html/capture/stream/dmitcat_track.php
user=www-data
autostart=false
autorestart=true
startretries=30
startsecs=4
